# docker/test/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; master_process on;"
autostart=true
autorestart=false
